# s6/acd_cli/finish
#!/bin/sh

/usr/bin/printf "Waiting for EncFS to unmount before unmounting ACD\n"
/bin/s6-svwait -d /etc/s6/encfs

/usr/bin/printf "Unmounting ACD\n"
/usr/bin/acd_cli umount
