# s6/encfs/finish
#!/bin/sh

/usr/bin/printf "unmounting EncFS\n"
/bin/fusermount -u /mnt/encrypted
