# s6/.s6-svscan/finish
#!/bin/sh

/bin/umount -f /mnt/encrypted
